# packages/bsp/hemai/wol@.service
[Unit]
Description=Wake-on-LAN for %i
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/sbin/ethtool -s %i wol g
Type=oneshot

[Install]
WantedBy=multi-user.target
